# ansible/roles/cluster-node/files/components/geard/systemd/container.slice

[Unit]
Description=Container slice container

[Slice]
CPUAccounting=yes
MemoryAccounting=yes
MemoryLimit=512M


[Install]
WantedBy=container.target container-active.target