# archived/pre-modules-pkg/systemd/container.slice

[Unit]
Description=Container slice container

[Slice]
CPUAccounting=yes
MemoryAccounting=yes
MemoryLimit=


[Install]
WantedBy=container.target container-active.target