# ansible/roles/cluster-node/files/components/nsq/systemd/nsqd.service
[Unit]

Description=Nsq daemon
Requires=vlipco.target
After=consul.service

[Service]

Restart=always
TimeoutStopSec=15
ExecStart=/bin/bash -c "/usr/local/bin/nsqd --broadcast-address=$(hostname).node.consul $(lookupd-params)"

[Install]

WantedBy=vlipco.target
