# package-src/components/base/misc/docker.service.d/dns_override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker -d -r=false --selinux-enabled \
  --dns 8.8.8.8 --dns 8.8.4.4 -H tcp://0.0.0.0:2375 -H unix:///var/run/docker.sock