# ansible/roles/cluster-node/files/components/common/systemd/docker.service.d/dns_override.conf
[Service]
ExecStart=
ExecStart=/usr/local/bin/docker -d -r=false --selinux-enabled \
  --dns 8.8.8.8 --dns 8.8.4.4