# archived/pre-modules-pkg/systemd/docker.service.d/dns_override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker -d -r=false --selinux-enabled \
  --dns 8.8.8.8 --dns 8.8.4.4