# package-src/core/etc/cluster.target
[Unit]
Description=Seaworthy Cluster Services
After=docker.service

[Install]
WantedBy=multi-user.target